# .config/systemd/user/wallpaper-changer.service
[Service]
Type=simple
ExecStart=/home/kani/.local/bin/wallpaper-changer --bing

[Install]
WantedBy=default.target
