# .config/systemd/wallpaper-changer.service
[Service]
Type=simple
ExecStart=/home/kani/.local/bin/wallpaper-changer --random

[Install]
WantedBy=default.target
